# supervisord.conf
[supervisord]
nodaemon=true

[program:postfix]
command=/opt/postfix.sh

[program:opendkim]
command=/usr/sbin/opendkim -f
